# systemd/doord.service
[Unit]
Description=Door daemon

[Service]
ExecStart=/usr/local/bin/doord

[Install]
WantedBy=multi-user.target
